# launchagents/com.missioncontrol.auto-spec.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.missioncontrol.auto-spec</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/natemacdaddy/.openclaw/workspace && set -a && source mission-control/.env.local && set +a && node mission-control/scripts/spec-from-deepseek.mjs</string>
  </array>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/natemacdaddy/.openclaw/workspace/logs/agents/auto-spec.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/natemacdaddy/.openclaw/workspace/logs/agents/auto-spec.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
